# app01_cpp/gateway-rpi/systemd/iot-gateway.service
# Plan Phase 5 (Deployment & Maintenance) packaging for the app01_cpp gateway.
#
# Install:
#   sudo cp gateway-rpi/build/gateway /usr/local/bin/iot-gateway
#   sudo cp gateway-rpi/systemd/iot-gateway.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now iot-gateway.service
#
# Dump per-node counters on demand (SR-1/SR-2 operator visibility, User Story 2):
#   sudo systemctl kill -s SIGUSR1 iot-gateway.service
#   journalctl -u iot-gateway.service -n 50
#
# Adjust ExecStart's serial device path for the deployed hardware.

[Unit]
Description=IoT sensor gateway (app01_cpp) - decodes and validates node telemetry
After=multi-user.target
# Cap the restart rate so a persistently failing unit doesn't spin forever.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/local/bin/iot-gateway --read /dev/ttyUSB0 --aggregate 12
# Full deployment (FR-3 publishing + Phase 5 monitoring + SR-6 auth):
#   ExecStart=/usr/local/bin/iot-gateway --read /dev/ttyUSB0 --aggregate 12 \
#       --metrics /var/lib/iot-gateway/gateway.prom \
#       --keys ${CREDENTIALS_DIRECTORY}/keys
#   StateDirectory=iot-gateway
#   LoadCredential=keys:/etc/iot-gateway/keys
# StateDirectory gives the sandboxed service its writable metrics directory;
# LoadCredential hands the root-owned 0600 key file to the DynamicUser
# process without loosening the key file's permissions.

# "Watchdog restart": systemd's own process supervision relaunches the
# gateway on crash or unexpected exit (rate-capped in [Unit]).
Restart=on-failure
RestartSec=2

# journald captures stdout/stderr automatically for systemd units; explicit
# here for documentation.
StandardOutput=journal
StandardError=journal

# Least privilege: no root. A transient unprivileged user, in the dialout
# group only for the serial device (Raspberry Pi OS puts /dev/ttyUSB* and
# /dev/ttyAMA* in group dialout).
DynamicUser=true
SupplementaryGroups=dialout

# Defense in depth (STRIDE: Elevation of Privilege): the process parses
# hostile bytes, so assume compromise and give it nothing — only the serial
# character device, no network, no capabilities, a native-arch service
# syscall allowlist, and read-only views of everything it doesn't own.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
DevicePolicy=closed
DeviceAllow=char-ttyUSB rw
DeviceAllow=char-ttyAMA rw
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
UMask=0077

[Install]
WantedBy=multi-user.target
